# projects/supportforge/supportforge.service
[Unit]
Description=SupportForge - Flat-rate Customer Support SaaS
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/support.beta.nex.monster/html
ExecStart=/usr/bin/node src/server.js
Restart=always
RestartSec=10
StandardOutput=append:/var/log/supportforge.log
StandardError=append:/var/log/supportforge.log
Environment=NODE_ENV=production
Environment=PORT=3847

[Install]
WantedBy=multi-user.target
